# packages/rocket-pi/rocket-pi-1.1.tar.gz/rocket-pi-1.1/rocketlauncher/rocket.service
# Project skidloader
# This is the systemd service description required for autostarting
# the rocketlaucher.sh at the boot time after the network services
# rocket.service requires to be located in etc/systemd/system

# Author: Martin Shishkov
# Created : 05.03.2020
# 
[Unit]
Description=Rocket launch service
After=multi-user.target

[Service]
Type=oneshot
ExecStart=bash %s
WorkingDirectory=%s
StandardOutput=inherit
StandardError=inherit
User=pi

[Install]
WantedBy=multi-user.target
